# archive/week2-revert-2026-03-05/launchd-disabled/com.miai.portfolio-health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.miai.portfolio-health</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/sheridanskala/.openclaw/workspace/scripts/run_service.sh</string>
        <string>/Users/sheridanskala/.openclaw/workspace/scripts/portfolio_health_monitor.py</string>
        
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/sheridanskala/.openclaw/workspace</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/sheridanskala/.openclaw/workspace/logs/com.miai.portfolio-health.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/sheridanskala/.openclaw/workspace/logs/com.miai.portfolio-health.err</string>
</dict>
</plist>
