# archive/week2-revert-2026-03-05/launchd-disabled/com.miai.factor-engine.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.miai.factor-engine</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/sheridanskala/.openclaw/workspace/scripts/factor_engine.py</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>StandardOutPath</key>
    <string>/Users/sheridanskala/.openclaw/workspace/logs/factor_engine.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/sheridanskala/.openclaw/workspace/logs/factor_engine_err.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
