# archive/week2-revert-2026-03-05/launchd-disabled/com.miai.intraday-momentum.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin</string>
	</dict>
	<key>Label</key>
	<string>com.miai.intraday-momentum</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/python3</string>
		<string>/Users/sheridanskala/.openclaw/workspace/scripts/intraday_momentum_scanner.py</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/sheridanskala/.openclaw/workspace/logs/momentum_scanner_err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/sheridanskala/.openclaw/workspace/logs/momentum_scanner.log</string>
	<key>StartInterval</key>
	<integer>600</integer>
	<key>WorkingDirectory</key>
	<string>/Users/sheridanskala/.openclaw/workspace</string>
</dict>
</plist>
