# crypto-trading-system/launchd/com.mi.crypto-stop-enforcer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mi.crypto-stop-enforcer</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/sheridanskala/.openclaw/workspace/crypto-trading-system/crypto_stop_enforcer.py</string>
        <string>--loop</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/sheridanskala/.openclaw/workspace/crypto-trading-system/logs/stop_enforcer_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/sheridanskala/.openclaw/workspace/crypto-trading-system/logs/stop_enforcer_stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/sheridanskala/.openclaw/workspace/crypto-trading-system</string>
</dict>
</plist>
